# deploy/systemd/jarvis-command-center.service
[Unit]
Description=JARVIS Master Command Center -- operator dashboard (stdlib HTTP server)
Documentation=file://__INSTALL_DIR__/deploy/HOST_RUNBOOK.md
After=network-online.target jarvis-live.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=-__INSTALL_DIR__/.env
Environment=PYTHONPATH=__PYTHONPATH__
Environment=JARVIS_MCC_HOST=127.0.0.1
Environment=JARVIS_MCC_PORT=8765
ExecStart=__INSTALL_DIR__/.venv/bin/python -m apex_predator.scripts.jarvis_dashboard
Restart=always
RestartSec=5
StandardOutput=append:%h/.local/log/apex_predator/jarvis-mcc.log
StandardError=append:%h/.local/log/apex_predator/jarvis-mcc.err

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/state/apex_predator %h/.local/log/apex_predator
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictNamespaces=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native

[Install]
WantedBy=default.target
